# runit/xui/run
#!/bin/sh
exec 2>&1
echo -e "starting x-ui..."
cd /usr/local/x-ui/ && exec /usr/local/x-ui/x-ui
