# runit/xui/run
#!/bin/sh
exec 2>&1
echo -e "starting x-ui..."
exec -c /usr/local/x-ui/x-ui
